# packaging/aur/librefang-docker/librefang-docker.service
[Unit]
Description=LibreFang Agent OS Container
Documentation=https://librefang.ai
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=-/etc/librefang/docker.env
ExecStart=/usr/bin/librefang-docker run
ExecStop=/usr/bin/librefang-docker stop
Restart=on-failure
RestartSec=5
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
